# config/crontab
* * * * * cd /var/www/drupal && export $(cat ../container.env | xargs) && vendor/bin/drush cron
